# dockerfiles/supervisor/redis.conf
[program:redis]
command=/usr/bin/redis-server 
autostart=true
autorestart=true