# dockerfiles/supervisor/restmq.conf
[program:restmq]
command=/srv/restmq/start_scripts/restmq_server --acl=/etc/restmq/acl.conf --listen=0.0.0.0
autostart=true
autorestart=true